# systemd/ruuvi.service
[Unit]
Description=RuuviTag metrics sender to InfluxDB
After=network-online.target bluetooth.service docker.service
Wants=network-online.target bluetooth.service docker.service

[Service]
Type=simple

User=<USER>

WorkingDirectory=/home/<USER>/raspi_server

ExecStart=/home/<USER>/raspi_server/venv/bin/python -u /home/<USER>/raspi_server/src/ruuvi_to_influx.py

Restart=always
RestartSec=5

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
